# infra/systemd/brh@.service
# Bridge Runtime Handler systemd service template
# 
# Installation:
# 1. Copy this file to /etc/systemd/system/brh@.service
# 2. Create /etc/brh/dominion.seal with your DOMINION_SEAL value
# 3. Enable and start the service:
#    sudo systemctl enable brh@"dominion://sovereign.bridge?env=prod&epoch=$(date +%s)&sig=<sig>"
#    sudo systemctl start brh@"dominion://..."

[Unit]
Description=Bridge Runtime Handler (%i)
After=docker.service
Requires=docker.service

[Service]
Environment=FORGE_DOMINION_ROOT=%i
Environment=DOMINION_SEAL=/etc/brh/dominion.seal
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/bridge-repo
ExecStart=/usr/bin/python3 -m brh.run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
