# bridge-catalog/temporal-anchor-20251207-172200Z/temporal-anchor-20251206-211650Z/temporal-anchor-20251206-210632Z/coefficient-20251208-0256-17r-3i-full-engine-startup/coefficient-20251208-0252-4r-3i-harmonic-ascent/coefficient-20251208-0246-9r-3i-final-resolution/coefficient-20251208-0241-31r-3i-full-system-audit-repair/coefficient-20251208-0232-4r-3i-self-description/coefficient-20251208-0214-28r-3i-development-weaving-fugue/coefficient-20251208-0155-23r-3i-subject-return-fugue/coefficient-20251208-0145-15r-3i-fourth-harmonic-overtone/coefficient-20251208-0130-18r-3i-third-harmonic-overtone/coefficient-20251208-0120-12r-3i-detailed-carving-leap/infra/systemd/brh@.service
# Bridge Runtime Handler systemd service template
# 
# Installation:
# 1. Copy this file to /etc/systemd/system/brh@.service
# 2. Create /etc/brh/dominion.env with:
#    DOMINION_SEAL=your-secret-seal-here
# 3. Enable and start the service:
#    sudo systemctl enable brh@"dominion://sovereign.bridge?env=prod&epoch=$(date +%s)&sig=<sig>"
#    sudo systemctl start brh@"dominion://..."

[Unit]
Description=Bridge Runtime Handler (%i)
After=docker.service
Requires=docker.service

[Service]
Environment=FORGE_DOMINION_ROOT=%i
EnvironmentFile=/etc/brh/dominion.env
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/bridge-repo
ExecStart=/usr/bin/python3 -m brh.run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
